# bengal/dot-config/systemd/user/mnt-claw.mount
[Unit]
Description=SSHFS mount for claw
After=network-online.target
Wants=network-online.target

[Mount]
What=claw:/home/claw
Where=/mnt/claw
Type=fuse.sshfs
Options=_netdev,user,delay_connect,reconnect,ServerAliveInterval=15,dir_cache=yes,idmap=user,follow_symlinks,transform_symlinks,compression=yes

[Install]
WantedBy=default.target
